# docs/install-file/ikuai-bypass.initd
#!/bin/sh /etc/rc.common

START=99
STOP=10

start() {
    /opt/ikuai-bypass/ikuai-bypass -r cronAft -c /opt/ikuai-bypass/config.yml > /dev/null 2>&1 &
}

stop() {
    killall -q -9 ikuai-bypass 2>/dev/null
}

restart() {
    stop
    sleep 1
    start
}
